# tenshistep-plasma/tools/systemd/tenshistep-update-notifier.service
[Unit]
Description=TenshiSTEP software-update notifier
Documentation=man:tenshistep-update-notifier

[Service]
Type=oneshot
ExecStart=%h/.local/bin/tenshistep-update-notifier --quiet
